# .config/systemd/user/swayidle.service
[Unit]
Description=Idle management daemon for Wayland
After=graphical-session.target
PartOf=graphical-session.target
Requisite=graphical-session.target

[Install]
WantedBy=niri.service

[Service]
ExecStart=/usr/bin/swayidle -w
Restart=on-failure
